# supervisord.conf
[supervisord]
nodaemon=true

[program:cupsd]
command=/usr/sbin/cupsd -f
priority=1
autostart=true
autorestart=true

[program:cupsctl]
command=bash -c "sleep 10 && /usr/sbin/cupsctl --remote-any ServerAlias=* DefaultEncryption=Never"
priority=2
autostart=true
autorestart=false
startsecs=0

[program:savapage]
command=bash -c "sleep 8; rm -f /opt/savapage/server/logs/service.pid; /opt/savapage/server/bin/linux-x64/app-server start; while pgrep -f org.savapage.server.WebServer >/dev/null; do sleep 10; done; echo 'savapage WebServer exited'; exit 1"
priority=3
autostart=true
autorestart=true
startsecs=30
user=savapage

[program:savapage-cups-notifier]
command=bash -c "sleep 15 && /opt/savapage/providers/cups/linux-x64/roottasks"
priority=4
autostart=true
autorestart=false
startsecs=0
user=root
